# crontab
* * * * * php /Users/cristianangulonova/Sites/Proyectos/elearn/app/console app:modulos:liberados
